# systemd/nvidia-persistenced.service.d/override.conf
# Drop-in override for nvidia-persistenced.service that makes it eGPU-aware.
#
# Problem: nvidia-persistenced is shipped as a system-wide daemon expected to
# run permanently — it starts at boot and fails noisily if /dev/nvidia* don't
# exist. On a laptop with an eGPU, the GPU isn't on the bus at boot. Without
# this drop-in:
#   - boot tries to start the service 5 times in a row
#   - each one fails with "Failed to query NVIDIA devices"
#   - systemd hits StartLimitBurst and gives up: "Start request repeated too
#     quickly"
#   - manual `systemctl reset-failed` is needed before any later attempt
#
# Fix: skip the service cleanly when no NVIDIA device file exists (boot path),
# and let the udev rule in udev/99-nvidia-egpu-persistenced.rules trigger the
# start/stop as the eGPU appears/disappears on the PCI bus.
#
# Install destination: /etc/systemd/system/nvidia-persistenced.service.d/override.conf

[Unit]
ConditionPathExists=/dev/nvidia0

[Service]
# udev handles re-starts on plug events; auto-restart on failure would just
# loop pointlessly when the GPU isn't there.
Restart=no
